# docker/supervisor.conf.d/janitoo_roomba.conf
[program:jnt_roomba]
user=root
autostart=true
autorestart=true
command=jnt_roomba -c /etc/janitoo/janitoo_roomba.conf front
