# docker/supervisor.conf.d/janitoo_roomba.conf
[program:jnt_roomba]
user=root
autostart=true
autorestart=true
command=jnt_roomba -c /opt/janitoo/etc/janitoo_roomba.conf start
